# 20_機構｜Mekhane/_src｜ソースコード/mekhane/deploy/systemd/hgk-ls-warmup-morning.timer
[Unit]
Description=HGK LS Warmup Morning 07:05 JST

[Timer]
OnCalendar=*-*-* 07:05:00
Persistent=true
AccuracySec=1m
Unit=hgk-ls-warmup.service

[Install]
WantedBy=timers.target
