# 20_機構｜Mekhane/_src｜ソースコード/mekhane/deploy/systemd/hgk-ls-warmup-evening.timer
[Unit]
Description=HGK LS Warmup Evening 17:05 JST

[Timer]
OnCalendar=*-*-* 17:05:00
Persistent=true
AccuracySec=1m
Unit=hgk-ls-warmup.service

[Install]
WantedBy=timers.target
