# 20_機構｜Mekhane/_src｜ソースコード/mekhane/deploy/systemd/hgk-ls-warmup-noon.timer
[Unit]
Description=HGK LS Warmup Noon 12:05 JST

[Timer]
OnCalendar=*-*-* 12:05:00
Persistent=true
AccuracySec=1m
Unit=hgk-ls-warmup.service

[Install]
WantedBy=timers.target
